# conf/celery.conf
[program:celery]
user=vagrant
command=celery worker -A cron --loglevel=info --beat
autostart=true
autorestart=true
directory=/home/vagrant/cron
