# conf/flower.conf
[program:flower]
user=vagrant
command=flower -A cron --port=5555
autostart=true
autorestart=true
directory=/home/vagrant/cron
